# scripts/launchd/com.copain.notion-pulse.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.copain.notion-pulse</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/FULL/PATH/TO/claude for me/scripts/notion-pulse.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/FULL/PATH/TO/claude for me</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/copain-notion-pulse.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/copain-notion-pulse.err</string>
</dict>
</plist>
